# router/pi-guardian-router.service
[Unit]
Description=PI Guardian Model Router
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=alex
Group=alex
WorkingDirectory=/home/alex/pi-guardian/router
EnvironmentFile=/home/alex/pi-guardian/router/.env
Environment=PYTHONPATH=/home/alex/pi-guardian/router
ExecStart=/bin/sh -c 'exec /home/alex/pi-guardian/router/.venv/bin/uvicorn app.main:app --host "${ROUTER_HOST:-127.0.0.1}" --port "${ROUTER_PORT:-8071}"'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
